# debian/hhvm.postrm
if [ "$1" = purge ]; then
   update-rc.d hhvm remove
fi

#DEBHELPER#
